# system_config/lib/80-import_usb.rules
KERNEL=="sd[b-z]1", ACTION=="add",ATTR{type}=="1", ATTR{timeout}="1800",  RUN+="/usr/local/sbin/pb_import.sh %k"
